# bin/easygsp
#!/bin/sh
#
# EasyGSP     Startup script for the groovy script server
#
# chkconfig: - 85 15
# description: Finally, the jvm in the LAMP stack
#
# processname: java
# config: /usr/local/share/easygsp/conf/server.properties


# Source function library
. /etc/rc.d/init.d/functions


prog="gss"
gss="/usr/local/share/easygsp/start.sh"
gssstop="/usr/local/share/easygsp/stop.sh"
RETVAL=0

start() {
	echo -n $"Starting $prog: "
	daemon $gss
	
	touch /var/lock/subsys/$prog
	return 
}

stop() {
	echo -n $"Stopping $prog: "
	$gssstop
	
        rm -f /var/lock/subsys/$prog
	return 
}


case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		start
		;;
	status)
		status $gss
		RETVAL=$?
		;;
	*)
		echo $"Usage: $0 {start|stop|restart|status}"
		RETVAL=1
esac

exit $RETVAL
